# files/boxx.service
[Unit]
Before=mysql.service

[Service]
ExecStart=/home/pressboxx/Config/boxx.sh

[Install]
WantedBy=default.target