# files/boxx.service
[Unit]
Before=mysql.service

[Service]
ExecStart=/home/boxx/Startup/boxx.sh

[Install]
WantedBy=default.target